# miner/bnc-miner.service
[Unit]
Description=Banncoin Miner (Canonical)
After=network.target

[Service]
Type=simple
User=banncoin
WorkingDirectory=/home/banncoin
ExecStart=/usr/bin/python3 /home/banncoin/bnc_miner.py \
  --api https://api.banncoin.org \
  --reward-to YOUR_WALLET_HERE \
  --log-dir /home/banncoin/miner_logs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
