# miner/bncminer.service
[Unit]
Description=Banncoin Miner (User-level variant)
After=network.target

[Service]
User=banncoin
WorkingDirectory=/home/banncoin
ExecStart=/usr/bin/python3 /home/banncoin/bnc_miner.py \
  --api https://api.banncoin.org \
  --reward-to YOUR_WALLET_HERE \
  --log-dir /home/banncoin/miner_logs
Restart=always
RestartSec=5
StandardOutput=append:/var/log/banncoin/miner.log
StandardError=append:/var/log/banncoin/miner.log

[Install]
WantedBy=multi-user.target
